# auto_trader/auto-trader.service
[Unit]
Description=Auto Trader - 국내 주식 자동매매 시스템
After=network.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/auto_trader
ExecStart=/home/trader/auto_trader/.venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
