# auto_trader/auto-trader.service
[Unit]
Description=Auto Trader - 국내 주식 자동매매 시스템
After=network.target

[Service]
Type=simple
User=tjrgn
WorkingDirectory=/home/tjrgn/workspace/auto-trade/auto_trader
ExecStart=/home/tjrgn/workspace/auto-trade/auto_trader/venv/bin/python3 main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
